# data/systemd-askpass/wayland-systemd-askpass.path
[Unit]
Description=Dispatch Password Requests to wayland-systemd-askpass

[Path]
DirectoryNotEmpty=%t/systemd/ask-password
MakeDirectory=yes
DirectoryMode=0750
